# data/crash-reporter-storagemon.path
[Unit]
Description=Crash files storage usage monitor
ConditionPathExists=!/tmp/os-update-running

[Path]
PathChanged=/var/cache/core-dumps

[Install]
WantedBy=paths.target
