# data/crash-reporter-endurance.service
[Unit]
Description=Crash reporter integration with sp-endurance

[Service]
Type=simple
ExecStart=/usr/libexec/endurance-collect
