# data/crash-reporter.service
[Unit]
Description=Crash file reporter daemon
Requires=dbus.service
After=user-session.target
ConditionPathExists=!/tmp/os-update-running

[Service]
ExecStart=/usr/bin/crash-reporter-daemon
Type=simple
Restart=always

[Install]
WantedBy=post-user-session.target
